# Consul/Agents/httpd/supervisor/confg.d/httpd.conf
[program:apache2]
command=httpd -DFOREGROUND

